# configs/supervisord/supervisor.d/nova-scheduler.ini
[program:nova-scheduler]
command=nova-scheduler
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nova-scheduler.log
