# configs/supervisord/supervisor.d/nova-compute.ini
[program:nova-compute]
command=nova-compute
