# configs/supervisord/supervisor.d/nova-conductor.ini
[program:nova-conductor]
command=nova-conductor
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nova-conductor.log
