# configs/supervisord/supervisor.d/nova-conductor.ini
[program:nova-conductor]
command=nova-conductor
