# configs/supervisord/supervisor.d/nova-compute.ini
[program:nova-compute]
command=nova-compute
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nova-compute.log

