# configs/supervisord/supervisor.d/nova-api-metadata.ini
[program:nova-api-metadata]
command=nova-api-metadata
