# configs/supervisord/supervisor.d/nova-spicehtml5proxy.ini
[program:nova-spicehtml5proxy]
command=nova-spicehtml5proxy
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nova-spicehtml5proxy.log
