# configs/supervisord/supervisor.d/nova-placement-api.ini
[program:nova-placement-api]
command=uwsgi --ini /etc/uwsgi/nova-placement-api.ini
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nova-placement-api.log

