# configs/supervisord/supervisor.d/nova-consoleauth.ini
[program:nova-consoleauth]
command=nova-consoleauth
