# configs/supervisord/supervisor.d/nova-cert.ini
[program:nova-cert]
command=nova-cert
redirect_stderr=true
stdout_logfile=/var/log/supervisord/nova-cert.log
