# configs/supervisord/supervisor.d/nova-scheduler.ini
[program:nova-scheduler]
command=nova-scheduler
