# configs/supervisord/supervisor.d/nova-spicehtml5proxy.ini
[program:nova-spicehtml5proxy]
command=nova-spicehtml5proxy

